# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn -c gunicorn_config.py app:app
directory=/app
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:keepalive]
command=sh -c "while true; do sleep 3600; done"
autostart=true
autorestart=false